# src/server/background-akonadi.slice
[Unit]
Description=Akonadi PIM Storage Service